# salt/salt/project/gunicorn.conf
[program:gunicorn]
command=/home/vagrant/.virtualenvs/{% templatetag openvariable %} projectname {% templatetag closevariable %}/bin/gunicorn {% templatetag openvariable %} projectname {% templatetag closevariable %}.wsgi:application --bind 127.0.0.1:1337 --workers 2 --timeout 30 --access-logfile /home/vagrant/{% templatetag openvariable %} projectroot {% templatetag closevariable %}/log/gunicorn-access.log --error-logfile /home/vagrant/{% templatetag openvariable %} projectroot {% templatetag closevariable %}/log/gunicorn-error.log --pythonpath /home/vagrant/{% templatetag openvariable %} projectroot {% templatetag closevariable %}/{% templatetag openvariable %} projectname {% templatetag closevariable %}/
directory=/home/vagrant/{% templatetag openvariable %} projectroot {% templatetag closevariable %}/{% templatetag openvariable %} projectname {% templatetag closevariable %}
user=vagrant
autostart=true
autorestart=true
redirect_stderr=True
environment=DJANGO_SETTINGS_MODULE="{% templatetag openvariable %} django_settings_module {% templatetag closevariable %}",SECRET_KEY="{% templatetag openvariable %} secret_key {% templatetag closevariable %}"
